# integration/robot-bridge/deploy/robot-bridge.service
[Unit]
Description=Outbound Robot Heartbeat Bridge
After=network-online.target
Wants=network-online.target

[Service]
User=robotbridge
WorkingDirectory=/opt/Intelligent-Power-Inspection/integration/robot-bridge
EnvironmentFile=/etc/robot-bridge/robot-bridge.env
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
